# deploy/alt-linux/systemd/alt-deploy-register.service
[Unit]
Description=ALT deployment registration API
After=network-online.target
Wants=network-online.target

Requires=alt-deploy-guard.service
After=alt-deploy-guard.service
[Service]
Type=simple
User=altserver
Group=altserver
Environment=PYTHONPATH=/opt/alt-deploy-control
ExecStart=/usr/bin/python3 /opt/alt-deploy-api/register_api.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/alt-deploy/registration
ReadWritePaths=/var/lib/alt-deploy

[Install]
WantedBy=multi-user.target
