# deploy/alt-linux/systemd/alt-deploy-process.path
[Unit]
Description=Watch for new ALT registrations

[Path]
DirectoryNotEmpty=/srv/alt-deploy/registration/pending
Unit=alt-deploy-process.service

[Install]
WantedBy=multi-user.target
