# deploy/alt-linux/systemd/alt-deploy-process.service
[Unit]
Description=Process pending ALT registrations
After=network-online.target alt-deploy-register.service
Wants=network-online.target

[Service]
Type=oneshot
User=altserver
Group=altserver
Environment=PYTHONPATH=/opt/alt-deploy-control
WorkingDirectory=/home/altserver/ansible
ExecStart=/usr/bin/python3 /opt/alt-deploy-api/process_pending.py
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/srv/alt-deploy/registration
ReadWritePaths=/home/altserver/.ssh
ReadWritePaths=/var/lib/alt-deploy
