# deploy/netctl-collect.service
[Unit]
Description=Network Observer collection

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/netctl --json collect all
