# deploy/netopsctl.service
[Unit]
Description=Netops control broker
Requires=netopsctl.socket
After=netopsctl.socket

[Service]
Type=simple
User=netopsctl
Group=netopsctl
SupplementaryGroups=netctl
EnvironmentFile=/etc/netopsctl/netopsctl.env
LoadCredential=netopsctl-audit-signing-key:/etc/netopsctl/credentials/audit_signing_ed25519.raw
WorkingDirectory=/opt/openvpn-web
ExecStart=/opt/openvpn-web/.venv/bin/python -m netopsctl.server
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/var/lib/netctl
ReadWritePaths=/var/lib/netopsctl /run/netopsctl
