# deploy/server-observer.service
[Unit]
Description=OpenVPN gateway server health observer
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=openvpm
Group=openvpn-web
WorkingDirectory=/
UMask=0027
TimeoutStartSec=3min
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=tmpfs
ProtectSystem=strict
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
CapabilityBoundingSet=
BindReadOnlyPaths=/etc/openvpn-web/server-observer.json
BindReadOnlyPaths=/etc/openvpn-web/server-observer.key
BindReadOnlyPaths=/etc/openvpn-web/server-observer.known_hosts
InaccessiblePaths=/etc/openvpn-web/openvpn-web.env
InaccessiblePaths=/etc/openvpn/client-generator
InaccessiblePaths=-/mnt/antares_soft/vpn_config
InaccessiblePaths=-/var/lib/openvpn-web/openvpn-web.sqlite
InaccessiblePaths=/opt/openvpn-web
InaccessiblePaths=-/var/lib/openvpn-web/server-drafts
InaccessiblePaths=/usr/local/lib/openvpn-web-server-draft-worker
ReadOnlyPaths=/usr/local/lib/openvpn-web-server-observer
ReadWritePaths=/var/lib/openvpn-web/server-observer
ExecStart=/usr/local/sbin/server-observer
