# deploy/vpn-policy.service
[Unit]
Description=Policy routing for VLAN50 through WireGuard
Wants=wg-quick@wg0.service
After=network-online.target wg-quick@wg0.service
PartOf=wg-quick@wg0.service

[Service]
Type=oneshot
ExecStart=/usr/bin/flock --exclusive /run/lock/vpn-policy.lock /usr/local/sbin/vpn-policy.sh start
ExecStop=/usr/bin/flock --exclusive /run/lock/vpn-policy.lock /usr/local/sbin/vpn-policy.sh stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
