# deploy/vpn-policy-reconcile.service
[Unit]
Description=Reconcile VLAN50 policy routing through WireGuard

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/sbin/vpn-policy.sh reconcile
StandardOutput=journal
StandardError=journal
